# cron/container_files/etc/cron.d/wayf-update
*/5 * * * * root /usr/local/bin/wayf-update.sh > /tmp/logcron 2>&1
